# terminal-logging/terminal-logs-rotate.timer
[Unit]
Description=Daily terminal-log rotation + R2 backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
